# systemd/bluelock.service
[Unit]
Description=bluelock BLE proximity lock/unlock daemon
After=bluetooth.target

[Service]
Type=simple
ExecStart=%h/.local/bin/bluelock run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
